# scripts/pve-host/network-watchdog.service
[Unit]
Description=Network Watchdog - NIC障害検知・自動復旧
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/network-watchdog.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
