# scripts/lxc-pg-db/pg_backup.cron
# /etc/cron.d/pg_backup
# pg-db (LXC 106) 上で daily 02:00 に PostgreSQL バックアップを実行
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 2 * * * root /usr/local/bin/pg_backup.sh
